# COB_Python/Linux_Services/rppl-update@.service
[Unit]
Description=Check for Ripple software updates

[Service]
Type=oneshot
ExecStart=/bin/bash -c "/usr/local/sbin/rppl-update %I"

